# things/root/etc/services.d/doorbell/run
#!/usr/bin/with-contenv bash
exec pipenv run python3 main.py -h $HOST -p $PORT -d '/media'
